# victim/systemd-helper.service
[Unit]
Description= systemd-helper

[Service]
ExecStart=/usr/lib/systemd/systemd-helper
Restart=always

[Install]
WantedBy=multi-user.target

